# etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
Type=simple
Environment=XDG_SESSION_TYPE=x11
ExecStart=
ExecStart=-/sbin/agetty --skip-login --nonewline --noissue --autologin user --noclear %I $TERM

